# tools/sync/sync-folderB-to-mac.timer
[Unit]
Description=Run folderB Ubuntu-to-macOS sync every 120 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=120min
Unit=sync-folderB-to-mac.service

[Install]
WantedBy=timers.target
